# deploy/systemd/netprospect-backfill-cms.service
[Unit]
Description=NetProspect — backfill CMS (pool principal, lotes de 50k controlados)
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/root/Github/netprospect-v1
EnvironmentFile=-/root/Github/netprospect-v1/docker/.env
Environment=NATS_URL=nats://localhost:4222
ExecStart=/root/Github/netprospect-v1/scripts/backfill-cms.sh
Restart=always
RestartSec=30
StandardOutput=append:/root/Github/netprospect-v1/out/backfill-cms.log
StandardError=append:/root/Github/netprospect-v1/out/backfill-cms.log

[Install]
WantedBy=default.target
