# deploy/observability/wordfence-update.service
[Unit]
Description=NetProspect — atualiza a Wordfence Intelligence vuln DB (índice compacto em MinIO)
# SEM After/Requires=docker.service: é uma unit de USER (systemctl --user) e o docker.service é de SISTEMA
# (não visível do systemd de user → falhava a arrancar). O `docker exec` funciona na mesma (daemon system-wide).

[Service]
Type=oneshot
# Corre DENTRO de um worker via o wrapper MOUNTADO worker/run-wordfence.mjs (worker/ é volume-mounted →
# sem rebuild; qualquer worker serve, tem MINIO_* + WORDFENCE_API_KEY + lib). No-op limpo se a key faltar.
ExecStart=/bin/sh -c 'c=$(docker ps --format "{{.Names}}" -f name=netprospect-worker- | head -1); test -n "$c" && docker exec "$c" node worker/run-wordfence.mjs'
